# .docker/php/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
minfds=10000

[program:php-fpm]
command = /usr/sbin/php-fpm8.3
autostart=true
autorestart=true

[program:messenger]
command=php /var/www/symfony/bin/console messenger:consume async --time-limit=3600 -vv
user=www-data
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
